# launchd/com.david.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.david.digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ohsay22/Developer/jobauto/venv/bin/python3</string>
        <string>/Users/ohsay22/Developer/jobauto/digest.py</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <!-- Run daily at 8am -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/ohsay22/Developer/jobauto</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/ohsay22</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ohsay22/anelo/logs/digest.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ohsay22/anelo/logs/digest.stderr.log</string>
</dict>
</plist>
